# test/integration/targets/systemd/templates/baz.service
[Unit]
Description=Baz Server
Documentation=Baz

[Service]
ExecStart=/bin/yes

[Install]
WantedBy=default.target
