# test/integration/targets/systemd/templates/dummy.socket
[Unit]
Description=Dummy Server Activation Socket

[Socket]
ListenDatagram=69

[Install]
WantedBy=sockets.target