# test/integration/targets/systemd/templates/mask_me.service
[Unit]
Description=Mask Me Server
Documentation=Mask

[Service]
ExecStart=/bin/yes

[Install]
WantedBy=default.target
