# test/integration/targets/systemd/templates/sleeper@.service
[Unit]
Description=Basic service to use as a template

[Service]
ExecStart={{ sleep_bin_path }} %i

[Install]
WantedBy=multi-user.target
