# test/integration/targets/service_facts/files/failed-test.service
[Unit]
Description=failed service test

[Service]
ExecStart=/usr/fakecmd

[Install]
WantedBy=multi-user.target