# test/integration/targets/systemd/files/dummy.service
[Unit]
Description=Dummy Service

[Service]
Type=simple
ExecStart=/bin/sleep 3600
