# test/integration/targets/systemd/files/sleeper@.service
[Unit]
Description=Basic service to use as a template

[Service]
ExecStart=/usr/bin/sleep %i

[Install]
WantedBy=multi-user.target
